# novnc-wp115/conf.d/115.conf
[program:115]
command=/bin/bash /usr/local/115/115.sh
autorestart=true
